# gcs-bkup/crontab
SHELL=/bin/bash
*/30 * * * * $HOME/backup.sh
